# supervisord.ini
[program:yourbot]
command=python3 footballOpenCV/main.py
directory= footballOpenCV/
autostart=true
autorestart=true
stderr_logfile=/var/log/yourbot/yourbot.err.log
stdout_logfile=/var/log/yourbot/yourbot.out.log